# bots/cron.conf

[program:cron]
command=/usr/sbin/service cron start

