# bots/node.conf


[program:npm]
command=/usr/local/bin/npm start

