# bots/svnserve.conf

[program:svnserve]
command=/usr/bin/svnserve --daemon --foreground --root /var/svn

