# node-selenium/fluxbox.conf

[program:fluxbox]
command=fluxbox -display :99

