# node-selenium/node.conf


[program:npm]
command=/usr/bin/npm start

