# node-selenium/x11vnc.conf

[program:x11vnc]
command=x11vnc -forever -usepw -shared -rfbport 5900 -display :99

